# systemd/logtriage-webui.service
[Unit]
Description=Log Triage Web UI
After=network-online.target logtriage.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/logtriage
Environment=LOGTRIAGE_CONFIG=/etc/logtriage/config.yaml
ExecStart=/opt/logtriage/venv/python3 -m logtriage.webui
Restart=on-failure
RestartSec=5s
StateDirectory=logtriage
RuntimeDirectory=logtriage
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
